# kitt/panp/panp.path
# copy to /usr/local/lib/systemd/system
# enable with systemctl enable panp.path

[Unit]
Description=KITT Power Auto Norm Pursuit Service

[Path]
PathExists=/dev/switchpod
PathExists=/dev/ttyAMA0
PathExists=/dev/ttyAMA1

[Install]
WantedBy=multi-user.target

